# mysql.conf
[supervisord]
nodaemon=true

[program:mysql]
priority=10
command=/usr/bin/mysqld_safe

[program:start]
priority=90
command=/start.sh
